# heatwave.service
[Unit]
Description=Heatwave MIDI+Flask App
After=network.target bluetooth.service
StartLimitIntervalSec=0

[Service]
Type=simple
User=heatwave
WorkingDirectory=/home/heatwave/heatwave-app
ExecStart=/home/heatwave/heatwave-venv/bin/python main.py
Restart=always
RestartSec=5
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=heatwave

[Install]
WantedBy=multi-user.target
